# infra/systemd/app.target
[Unit]
Description=app
Requires=app.service
