# infra/systemd/secrets.target
[Unit]
Description=secrets
Requires=secrets.service

[Install]
WantedBy=multi-user.target
