# infra/systemd/worker.target
[Unit]
Description=worker
Requires=worker.service
