# infra/systemd/app.service.d/production.conf
[Service]
Environment="DATABASE_URL=postgres://westrikworld_user:uifJ7TTAT3iFxEmnCmBA@localhost/westrikworld_app"
Environment="PASSWORD_HASH_SALT=z5oKODAWvspQrOWre8YB"
